# boetheia.service
[Unit]
Description=Boetheia

[Service]
User=pi
ExecStart=/home/pi/boetheia/boetheia.py

[Install]
WantedBy=default.target

